# weekly_review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tloughran.summa-weekly-review</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/tomloughran/Documents/Claude/Projects/RC Karpathy Wiki Project/scripts/generate_weekly_review.py</string>
  </array>

  <!-- Runs every Sunday at 20:00 local time — one hour before the nightly sync -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key> <integer>0</integer>
    <key>Hour</key>    <integer>20</integer>
    <key>Minute</key>  <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/tomloughran/Documents/Claude/Projects/RC Karpathy Wiki Project/weekly_review.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tomloughran/Documents/Claude/Projects/RC Karpathy Wiki Project/weekly_review.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
